# Hacking/Club/ToroHack_8/Docker_CTF/flask-api.service
[Unit]
Description=Flask API Service
After=network.target

[Service]
User=root
WorkingDirectory=/app
ExecStart=/usr/bin/python3 /app/api.py
Restart=always

[Install]
WantedBy=multi-user.target

